# crontab/jenkins.crontab
# m h  dom mon dow   command
1  20  *  *  *  jenkinsBackup.sh 2>&1 > /mnt/backupJenkins.log
45  21  *  *  * s3sync.sh 2>&1 > /mnt/sync.log
